# agents/roomler-agent/packaging/linux/roomler.service
[Unit]
Description=Roomler AI remote-control agent
Documentation=https://roomler.ai/
# An interactive user session is required for scrap capture (XShm needs a
# display) and enigo input injection (both run on the user's X server).
# This is a systemd *user* unit, not a system unit, for that reason.
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
# P3d Slice B: the daemon binary is now `roomlerd` (was `roomler-agent`).
ExecStart=/usr/bin/roomlerd run
Restart=on-failure
RestartSec=5
Environment=RUST_LOG=roomler_agent=info,warn

# Sandboxing — conservative defaults appropriate for a user unit.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
# P3d Slice B: the appdirs resolver (src/appdirs.rs) uses the NEW `roomler`
# app segment on fresh installs but keeps the OLD `roomler-agent` tree when a
# pre-rename host already has one. Whitelist BOTH so neither a fresh install
# nor an upgraded host is blocked by ProtectHome=read-only.
ReadWritePaths=%h/.config/roomler %h/.local/share/roomler %h/.config/roomler-agent %h/.local/share/roomler-agent

[Install]
WantedBy=default.target
